# ex07/96-ex07.rules
ACTION=="add", DEVPATH=="/module/debugfs", SUBSYSTEM=="module", RUN+="/root/little_penguin_1/ex07/udev.sh"
